# docker/monolith/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
priority=10
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data
user=postgres
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql.err

[program:redis]
priority=20
command=redis-server --daemonize no --bind 0.0.0.0
user=redis
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis.err

[program:api]
directory=/app
command=node apps/api/dist/main
user=root
autorestart=true
startretries=5
startsecs=5
environment=NODE_ENV=production,NODE_PATH="/app/node_modules/.pnpm/node_modules:/app/node_modules",DATABASE_URL="postgresql://jasca:jasca_secret@127.0.0.1:5432/jasca?sslmode=disable",REDIS_URL="redis://127.0.0.1:6379",JWT_SECRET="jasca_offline_secret",CORS_ORIGIN="https://jasca.koreacb.com",TRIVY_CACHE_DIR="/app/trivy-db"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
directory=/app
command=node server.js
user=root
autorestart=true
startretries=5
startsecs=10
; API_URL is used for server-side requests inside the container.
environment=NODE_ENV=production,NODE_PATH="/app/node_modules/.pnpm/node_modules:/app/node_modules",PORT=3000,HOSTNAME="0.0.0.0",API_URL="http://127.0.0.1:3001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
